# crontab
@reboot node /root/flycapture/newserver.js
